# etc/cron.d/tailscale-kill-switch
*/15 * * * * root /usr/local/bin/tailscale-kill-switch update
